# .config/systemd/user/accent-picker.service
[Unit]
Description=macOS-style press-and-hold accent picker
After=default.target

[Service]
ExecStart=%h/scripts/fedora/accent-picker.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
