# .config/systemd/system-units/btrfs-scrub.service
[Unit]
Description=Btrfs scrub on /
[Service]
Type=oneshot
ExecStart=/usr/sbin/btrfs scrub start -B /
