# .config/systemd/system-units/btrfs-scrub.timer
[Unit]
Description=Monthly Btrfs scrub on /

[Timer]
OnCalendar=monthly
Persistent=true
RandomizedDelaySec=1d

[Install]
WantedBy=timers.target
